# files/autostart/org.ferdium.Ferdium.desktop
[Desktop Entry]
Name=Ferdium
Exec=flatpak run org.ferdium.Ferdium
Terminal=false
Type=Application
Icon=org.ferdium.Ferdium